# echostream.service
[Unit]
Description=EchoStream Audio Communication
After=network.target mosquitto.service sound.target
Wants=mosquitto.service

[Service]
Type=simple
User=naveen
WorkingDirectory=/home/naveen/Linux_EchoStream
ExecStart=/home/naveen/Linux_EchoStream/api_call
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal

# GPIO permissions
SupplementaryGroups=gpio

# Audio permissions
SupplementaryGroups=audio

# Environment variables
Environment=DISPLAY=:0

[Install]
WantedBy=multi-user.target
